# src/install/files/conductor-service@.service
[Unit]
Description=Conductor service %f
Wants=conductor-service-config-config@%f.service

[Service]
WorkingDirectory=%f
TimeoutStartSec=0
Restart=always
Type=notify-reload
NotifyAccess=all
TimeoutSec=300
# RemainAfterExit=yes

ExecStart=/usr/local/bin/conductor _ service start
ExecReload=/usr/local/bin/conductor _ service restart
ExecStop=/usr/local/bin/conductor _ service stop
ExecStopPost=/usr/local/bin/conductor _ service cleanup

[Install]
WantedBy=multi-user.target

