# src/install/files/conductor-policy-server.socket
[Socket]
ListenStream=%t/conductor-policy.socket

[Install]
WantedBy=sockets.target
